# mkosi.profiles/bootc/mkosi.extra/usr/lib/systemd/system/systemd-homed-firstboot.service.d/make-admin.conf
[Service]
ExecStart=
ExecStart=homectl firstboot --prompt-new-user --prompt-shell=no --prompt-groups=no --mute-console=yes -G wheel
